# scripts/systemd/info-integrated.service
[Unit]
Description=AIBackgroundWorker Integrated Pipeline (Analyze + Deep + Report)
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/perso/analysis/AIBackgroundWorker
Environment="PATH=/home/perso/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/home/perso/analysis/AIBackgroundWorker/scripts/info_collector/integrated_pipeline.sh --analyze-batch-size 30 --deep-batch-size 3 --deep-min-importance 0.0 --deep-min-relevance 0.0 --report-hours 1

[Install]
WantedBy=multi-user.target
